# server/supervisor.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
loglevel=debug
logfile_maxbytes=0

[program:emitter-test-server]
command=emitter
autostart=true
autorestart=true
redirect_stderr=true
